# www-apps/evebox-bin/files/evebox
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="EveBox Server"
command=/usr/bin/evebox
command_args="server -c /etc/evebox/evebox.yaml"
command_background="true"
pidfile="/run/${RC_SVCNAME}.pid"
export EVEBOX_DATA_DIRECTORY=/var/lib/evebox

depend() {
	need net
}